# server_side/dockerfiles/memoryoftheworld/nginx/nginx.conf
[program:set_prosody_ip]
command=/bin/sed -i 's/PROSODY_PORT_5281_TCP_ADDR/'`echo $PROSODY_PORT_5281_TCP_ADDR`'/' /etc/nginx/sites-enabled/bosh
priority=2
startsecs=0
exitcodes=0

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_error.log
autostart=true
autorestart=true

